# remap/xremap.service
[Service]
ExecStart=/home/chs/dotfiles/remap/xremap /home/chs/dotfiles/remap/xremap.config
Restart=always

[Install]
WantedBy=gdm.service
